# data/fedora_update_utility.service
[Unit]
Description=Fedora Update Utility

[Service]
Type=oneshot
ExecStart=@exec@ --daemon

[Install]
WantedBy=graphical.target