# data/fedora_update_utility.service
[Unit]
Description=Fedora Update Utility

[Service]
ExecStart=/usr/bin/fedora_update_utilty

[Install]
WantedBy=default.target